# deploy/systemd/source-health-snapshot.service
[Unit]
Description=Finance Agent daily data source health snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/workspace/finance-agent
Environment=no_proxy=127.0.0.1,localhost,::1
Environment=UV_CACHE_DIR=/tmp/uv-cache
ExecStart=/usr/bin/env uv run python scripts/record_data_source_health_snapshot.py

[Install]
WantedBy=default.target
